# deploy/jazz2.service
[Unit]
Description=Jazz2 Next.js (cafe site)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/jazz2-app
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/node /var/www/jazz2-app/node_modules/next/dist/bin/next start -H 127.0.0.1 -p 3001
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
